# basic/common/files/etc/profile.d/bash_history.sh
#!/bin/bash

export HISTTIMEFORMAT='[%h %d %H:%M:%S] '
