# basic/git-autocommit/templates/etc/cron.d/nxs-git-autocommit.j2
* *	* * *	root	/usr/bin/dpkg -l > /etc/package_list
* *	* * *	root	sleep 10;cd / && git add -A && git commit -m "Autocommit @system" > /dev/null

